# deploy/hvac-backend.service
[Unit]
Description=944S HVAC backend (control loop + WebSocket + dashboard server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mark
WorkingDirectory=/home/mark/hvac
ExecStart=/usr/bin/python3 /home/mark/hvac/hvac_backend.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
